# systemd/hover-ddns.timer
[Unit]
Description=Timer for Hover DDNS Updater

[Timer]
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target

